# github-mirror/systemd/github-mirror-fallback.service
[Unit]
Description=ATS CI GitHub mirror fallback refresh
Documentation=file:@INSTALL_ROOT@/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@INSTALL_ROOT@
ExecStart=/usr/bin/docker-compose run --rm github-mirror-fallback
TimeoutStartSec=30min
Nice=5
