# github-mirror/systemd/github-mirror-webhook.service
[Unit]
Description=ATS CI GitHub mirror webhook receiver
Documentation=file:@INSTALL_ROOT@/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@MIRROR_USER@
Group=@MIRROR_GROUP@
EnvironmentFile=/etc/trafficserver-github-mirror/github-mirror-webhook.env
WorkingDirectory=@MIRROR_ROOT@
ExecStart=@INSTALL_ROOT@/bin/github-mirror-webhook.py
Restart=on-failure
RestartSec=5s
RestartPreventExitStatus=1
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=@MIRROR_ROOT@

[Install]
WantedBy=multi-user.target
